# rules/80-hotplug-cpu-mem.rules
# do not edit this file, it will be overwritten on update

# Hotplug physical CPU
SUBSYSTEM=="cpu", ACTION=="add", TEST=="online", ATTR{online}=="0", \
 ATTR{online}="1"

# Hotplug physical memory
SUBSYSTEM=="memory", ACTION=="add", TEST=="state", ATTR{state}=="offline", \
 ATTR{state}="online", TAG+="tmpfs"

#
TAG=="tmpfs", RUN+="/usr/lib/udev/remount-tmpfs"
